# rag-web.service
[Unit]
Description=RAG Web Application
After=network.target

[Service]
Type=simple
User=jarvis
WorkingDirectory=/home/jarvis/rag
EnvironmentFile=/home/jarvis/rag/.env
ExecStart=/home/jarvis/rag/venv/bin/gunicorn \
    --workers 1 \
    --threads 4 \
    --timeout 120 \
    --bind 0.0.0.0:8080 \
    --access-logfile - \
    --error-logfile - \
    wsgi:application
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
